# services/yum-sync.service
[Unit]
Description=Sync YUM Repo from rclone remote

[Service]
Type=oneshot
ExecStart=/home/hansel/scripts/sync-yum.sh
